# scripts/com.arec.morningbriefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arec.morningbriefing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>PROJECT_DIR_PLACEHOLDER/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PROJECT_DIR_PLACEHOLDER</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>PROJECT_DIR_PLACEHOLDER</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/oscar/Library/Logs/arec-morning-briefing.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/oscar/Library/Logs/arec-morning-briefing-error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
